# src/skel/etc/supervisor/conf.d/program-dropbear.conf
# -*- ini -*-
[program:dropbear]
command = /usr/local/bin/launch-dropbear.sh
startretries = 100000000
autorestart = true

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
