# src/skel/etc/supervisor/conf.d/program-darkice.conf
# -*- ini -*-
[program:darkice]
command = /usr/local/bin/launch-darkice.sh
startretries = 100000000
autorestart = true
