# src/skel/etc/supervisor/conf.d/program-dropbear.conf
# -*- ini -*-
[program:dropbear]
command = /usr/local/bin/launch-dropbear.sh
startretries = 100000000
autorestart = true
